# opensecuritycam.service
[Unit]
Description=OpenSecurityCam motion detection service
After=network.target remote-fs.target
Requires=remote-fs.target

[Service]
Type=simple
User=motion
Group=motion
Environment=UV_PYTHON_INSTALL_DIR=/usr/local/share/uv-python
# On kernel 6.12 (Debian Trixie) the unicam driver does not propagate the
# subdev format to the video node automatically — the node resets to Y12 on
# each boot.  motion will fail VIDIOC_STREAMON unless the correct packed
# 10-bit Bayer format (pGAA / SGBRG10P) is set first.  The '+' prefix runs
# this command as root regardless of User= above.
ExecStartPre=+/usr/bin/v4l2-ctl -d /dev/video0 --set-fmt-video=width=640,height=480,pixelformat=pGAA
ExecStart=/usr/bin/motion -c /etc/cam_motion/motion.conf -n
Restart=on-failure
RestartSec=5
StandardOutput=null
StandardError=journal

[Install]
WantedBy=multi-user.target
